# scripts/hypatia-autofix.timer
# SPDX-License-Identifier: PMPL-1.0-or-later
# Runs hypatia auto-fix weekly (Sunday 02:00)
[Unit]
Description=Weekly Hypatia auto-fix scan

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
